# edx/app/supervisor/conf.available.d/workers.conf
[program:cms_default_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=production,LANG=en_US.UTF-8,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=cms,BOTO_CONFIG="/edx/app/edxapp/.boto,EDX_REST_API_CLIENT_NAME=edx.cms.core.default,"
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/worker.sh cms --settings=production celery worker --loglevel=info --queues=edx.cms.core.default --hostname=edx.cms.core.default.%%h --concurrency=1 
killasgroup=true
stopwaitsecs=432000
; Set autorestart to `true`. The default value for autorestart is `unexpected`, but celery < 4.x will exit
; with an exit code of zero for certain types of unrecoverable errors, so we must make sure that the workers
; are auto restarted even when exiting with code 0.
; The Celery bug was reported in https://github.com/celery/celery/issues/2024, and is fixed in Celery 4.0.0.
autorestart=true

[program:cms_high_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=production,LANG=en_US.UTF-8,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=cms,BOTO_CONFIG="/edx/app/edxapp/.boto,EDX_REST_API_CLIENT_NAME=edx.cms.core.high,"
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/worker.sh cms --settings=production celery worker --loglevel=info --queues=edx.cms.core.high --hostname=edx.cms.core.high.%%h --concurrency=1 
killasgroup=true
stopwaitsecs=432000
; Set autorestart to `true`. The default value for autorestart is `unexpected`, but celery < 4.x will exit
; with an exit code of zero for certain types of unrecoverable errors, so we must make sure that the workers
; are auto restarted even when exiting with code 0.
; The Celery bug was reported in https://github.com/celery/celery/issues/2024, and is fixed in Celery 4.0.0.
autorestart=true

[program:lms_default_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=production,LANG=en_US.UTF-8,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms,BOTO_CONFIG="/edx/app/edxapp/.boto,EDX_REST_API_CLIENT_NAME=edx.lms.core.default,"
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/worker.sh lms --settings=production celery worker --loglevel=info --queues=edx.lms.core.default --hostname=edx.lms.core.default.%%h --concurrency=1 
killasgroup=true
stopwaitsecs=432000
; Set autorestart to `true`. The default value for autorestart is `unexpected`, but celery < 4.x will exit
; with an exit code of zero for certain types of unrecoverable errors, so we must make sure that the workers
; are auto restarted even when exiting with code 0.
; The Celery bug was reported in https://github.com/celery/celery/issues/2024, and is fixed in Celery 4.0.0.
autorestart=true

[program:lms_high_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=production,LANG=en_US.UTF-8,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms,BOTO_CONFIG="/edx/app/edxapp/.boto,EDX_REST_API_CLIENT_NAME=edx.lms.core.high,"
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/worker.sh lms --settings=production celery worker --loglevel=info --queues=edx.lms.core.high --hostname=edx.lms.core.high.%%h --concurrency=1 
killasgroup=true
stopwaitsecs=432000
; Set autorestart to `true`. The default value for autorestart is `unexpected`, but celery < 4.x will exit
; with an exit code of zero for certain types of unrecoverable errors, so we must make sure that the workers
; are auto restarted even when exiting with code 0.
; The Celery bug was reported in https://github.com/celery/celery/issues/2024, and is fixed in Celery 4.0.0.
autorestart=true

[program:lms_high_mem_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=production,LANG=en_US.UTF-8,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms,BOTO_CONFIG="/edx/app/edxapp/.boto,EDX_REST_API_CLIENT_NAME=edx.lms.core.high_mem,"
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/worker.sh lms --settings=production celery worker --loglevel=info --queues=edx.lms.core.high_mem --hostname=edx.lms.core.high_mem.%%h --concurrency=1 --maxtasksperchild 1
killasgroup=true
stopwaitsecs=432000
; Set autorestart to `true`. The default value for autorestart is `unexpected`, but celery < 4.x will exit
; with an exit code of zero for certain types of unrecoverable errors, so we must make sure that the workers
; are auto restarted even when exiting with code 0.
; The Celery bug was reported in https://github.com/celery/celery/issues/2024, and is fixed in Celery 4.0.0.
autorestart=true


[group:edxapp_worker]
programs=cms_default_1,cms_high_1,lms_default_1,lms_high_1,lms_high_mem_1
